# sysext/root/usr/lib/systemd/system/containerd.service.d/10-use-custom-config.conf
[Service]
ExecStartPre=cp /usr/share/containerd/config.toml /etc/containerd.toml
ExecStartPre=/usr/bin/bash -c "cat /usr/share/cloudhv/etc/containerd/config.toml.snippet >> /etc/containerd.toml"
Environment=CONTAINERD_CONFIG=/etc/containerd.toml
